# contrib/init/epicd.init
#!/bin/bash
#
#  epicd The Epic core server.
#
#
# chkconfig: 345 80 20
# description: epicd
# processname: epicd
#

# Source function library.
. /etc/init.d/functions

# you can override defaults in /etc/sysconfig/epicd, see below
if [ -f /etc/sysconfig/epicd ]; then
        . /etc/sysconfig/epicd
fi

RETVAL=0

prog=epicd
# you can override the lockfile via BITCOIND_LOCKFILE in /etc/sysconfig/epicd
lockfile=${BITCOIND_LOCKFILE-/var/lock/subsys/epicd}

# epicd defaults to /usr/bin/epicd, override with BITCOIND_BIN
bitcoind=${BITCOIND_BIN-/usr/bin/epicd}

# epicd opts default to -disablewallet, override with BITCOIND_OPTS
bitcoind_opts=${BITCOIND_OPTS}

start() {
    echo -n $"Starting $prog: "
    daemon $DAEMONOPTS $bitcoind $bitcoind_opts
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && touch $lockfile
    return $RETVAL
}

stop() {
    echo -n $"Stopping $prog: "
    killproc $prog
    RETVAL=$?
    echo
    [ $RETVAL -eq 0 ] && rm -f $lockfile
    return $RETVAL
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        status $prog
        ;;
    restart)
        stop
        start
        ;;
    *)
        echo "Usage: service $prog {start|stop|status|restart}"
        exit 1
        ;;
esac
